# claude-cli-bridge/claude-bridge.service
[Unit]
Description=Claude CLI Bridge Service
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=jay
WorkingDirectory=/home/jay/claude-cli-bridge
Environment=CLAUDE_CLI_PATH=/usr/local/bin/claude
ExecStart=/usr/bin/python3 /home/jay/claude-cli-bridge/claude_bridge.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
